# Desktop/workspace/bsf-looptech/config/com.bsf-looptech.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- BSF-LoopTech システム監視（30分ごと） -->

    <key>Label</key>
    <string>com.bsf-looptech.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tonton/Desktop/workspace/bsf-looptech/scripts/system_monitor.sh</string>
        <string>--quiet</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tonton/Desktop/workspace/bsf-looptech</string>

    <!-- 30分ごとに実行 (1800秒) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- ログ設定 -->
    <key>StandardOutPath</key>
    <string>/Users/tonton/Desktop/workspace/bsf-looptech/logs/monitor_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tonton/Desktop/workspace/bsf-looptech/logs/monitor_stderr.log</string>

    <!-- 環境変数 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>BSF_BACKUP_DIR</key>
        <string>/Users/tonton/BSF_Backups</string>
    </dict>

    <key>UserName</key>
    <string>tonton</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
